# deploy/systemd/op-chat-server.service
[Unit]
Description=OP-DBUS Chat Server
Documentation=https://github.com/ghostbridgetech/op-dbus
After=network.target dbus.service
Wants=dbus.service

[Service]
Type=simple
User=op-dbus
Group=op-dbus
EnvironmentFile=-/etc/op-dbus/environment
EnvironmentFile=-/etc/op-dbus/secrets.env
WorkingDirectory=/usr/local/sbin
ExecStart=/usr/local/sbin/op-web-server
Restart=always
RestartSec=5
WatchdogSec=30

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

# Allow network and D-Bus
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service

# Read-write paths
ReadWritePaths=/var/lib/op-dbus /var/log/op-dbus /run/op-dbus

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=op-chat-server

[Install]
WantedBy=multi-user.target
