# deploy/systemd/op-mcp.service
[Unit]
Description=OP MCP Aggregator (HTTP + Chat Control Bridge)
Documentation=file:/home/jeremy/op-dbus-v2/ANTIGRAVITY-MCP-README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/op-dbus
Environment=RUST_LOG=info
ExecStart=/usr/local/sbin/op-mcp-server --sse --bind 0.0.0.0:8083
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
