# deploy/systemd/op-web-servers.service
[Unit]
Description=op-web Servers Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jeremy
Group=jeremy
WorkingDirectory=/home/jeremy/op-dbus-v2

# Environment
Environment="PORT=8084"
Environment="STATIC_DIR=/home/jeremy/op-dbus-v2/static"
Environment="RUST_LOG=info"
EnvironmentFile=/home/jeremy/.op-web.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jeremy/op-dbus-v2
CapabilityBoundingSet=

# Process limits
LimitNOFILE=65536
LimitNPROC=512

# Start service
ExecStart=/home/jeremy/op-dbus-v2/target/release/op-web-server

# Restart policy
Restart=always
RestartSec=10s
StartLimitInterval=5min
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=op-web-servers

[Install]
WantedBy=multi-user.target
