# deploy/systemd/op-web.service
[Unit]
Description=op-dbus-v2 Web Server
Documentation=https://github.com/your-repo/op-dbus-v2
After=network-online.target dbus.service nginx.service
Wants=network-online.target dbus.service
Requires=dbus.socket

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/op-dbus
EnvironmentFile=/etc/op-dbus/environment

# Main process
ExecStart=/usr/local/sbin/op-web-server

# Restart policy
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=false
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
# RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
SystemCallArchitectures=native

# Read-write paths
ReadWritePaths=/var/cache/op-dbus
ReadWritePaths=/tmp

# Resource limits
LimitNOFILE=65535
LimitNPROC=512
MemoryMax=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=op-web

# Watchdog disabled - LLM requests can take >60s with multi-step tool execution
# WatchdogSec=60

[Install]
WantedBy=multi-user.target
